# xfce4/.config/autostart/resfix.desktop
[Desktop Entry]

Type=Application
Name=Screen Resolution Fix
Exec=$HOME/.bin/resfix
Terminal=false