# xfce4/.config/autostart/zentile.desktop
[Desktop Entry]

Type=Application
Name=Zentile
Exec=zentile
Terminal=false